# EXCEL_SEARCH_NON_TABULAR/excel-monitor.service
[Unit]
Description=Excel File Monitor and Database Uploader
After=network.target postgresql.service
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple
User=your_username
Group=your_group
WorkingDirectory=/home/omar/Desktop/Scripts/EXCEL_SEARCH_NON_TABULAR
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/omar/Desktop/Scripts/EXCEL_SEARCH_NON_TABULAR

# Execution
ExecStart=/home/omar/Desktop/Scripts/EXCEL_SEARCH_NON_TABULAR/venv/bin/python folder_monitor.py
ExecReload=/bin/kill -HUP $MAINPID

# Restart configuration
Restart=on-failure
RestartSec=5s
RestartPreventExitStatus=255
TimeoutStartSec=0

# Recovery and persistence
RuntimeDirectory=excel_monitor
RuntimeDirectoryMode=0755
StateDirectory=excel_monitor
StateDirectoryMode=0750
CacheDirectory=excel_monitor
CacheDirectoryMode=0750
LogsDirectory=excel_monitor
LogsDirectoryMode=0755

# Resource limits
LimitNOFILE=65535
TasksMax=infinity
TimeoutStopSec=5min

# Security settings
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes

# Logging
StandardOutput=append:/var/log/excel_monitor/service.log
StandardError=append:/var/log/excel_monitor/error.log

[Install]
WantedBy=multi-user.target